# external_loadbalancer_hap/system/hap_manager.service
#
# Systemd unit file for HAP Manager
#

[Unit]
Description=HAP Manager

[Service]
<<<<<<< HEAD
User=root
=======
>>>>>>> 324fe55465cb37533e0d50a571dbc2358458b4b4
Type=forking
ExecStart=/usr/bin/python /Celtic-Kubernetes/external_loadbalancer_hap/hap_manager_daemon.py start
ExecStop=/usr/bin/python /Celtic-Kubernetes/external_loadbalancer_hap/hap_manager_daemon.py stop
Restart=/usr/bin/python /Celtic-Kubernetes/external_loadbalancer_hap/hap_manager_daemon.py stop; /usr/bin/python /Celtic-Kubernetes/external_loadbalancer_hap/hap_manager_daemon.py start
PIDFile=/var/run/hap_manager.pid

[Install]
WantedBy=multi-user.target